# cloudflare-dns-bot.service
[Unit]
Description=Cloudflare DNS Manager Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/workspace
ExecStart=/workspace/venv/bin/python /workspace/Main.py
Restart=always
RestartSec=5

Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target